# root-basic/etc/inittab
#######################################
#                                     #
#    term:runlevels:action:command    #
#                                     #
#######################################

console::sysinit:-/etc/init.d/rcS
console::respawn:/sbin/getty -L ttyS1 115200 vt102

